# demo/ipython.supervisord.conf.j2
[program:ipython]
command=/home/ipython/bin/ipython notebook
autorestart=true
user=ipython
directory=/home/ipython
